# deploy/systemd/vps-sentry.service
# /etc/systemd/system/vps-sentry.service
[Unit]
Description=VPS Sentry (security monitor)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

User=root
Group=root
Nice=10
UMask=0077

# State dir (/var/lib/vps-sentry) created automatically
StateDirectory=vps-sentry
StateDirectoryMode=0700

# Avoid writing __pycache__ under /home (we keep ProtectSystem=strict)
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1

# IMPORTANT: don't let anything touch /root/.ssh (known_hosts, etc)
# Put HOME somewhere neutral that we already allow writing.
Environment=HOME=/var/lib/vps-sentry
Environment=SSH_AUTH_SOCK=

# Run script (installed wrapper)
ExecStart=/usr/local/bin/vps-sentry --format text

# ---- Hardening (safe defaults for your use case) ----
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes

ProtectSystem=strict
ProtectHome=read-only

ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes

LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictRealtime=yes
SystemCallArchitectures=native

# Allow only the families you actually need (python + netlink + inet + unix)
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

# Be explicit about what can be written (StateDirectory already does this, but this is belt+suspenders)
ReadWritePaths=/var/lib/vps-sentry

# Keep logs from going insane if something loops
LogRateLimitIntervalSec=30s
LogRateLimitBurst=200

# Guard rails
TimeoutStartSec=45s
