# deploy/systemd/vps-sentry-unit-event@.service
[Unit]
Description=VPS Sentry unit event hook for %i
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
Nice=10
UMask=0077

StateDirectory=vps-sentry
StateDirectoryMode=0711
RuntimeDirectory=vps-sentry
RuntimeDirectoryMode=0755

Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/var/lib/vps-sentry
Environment=SSH_AUTH_SOCK=

ExecStart=/usr/local/bin/vps-sentry-service-event --spec %i

NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes

ProtectSystem=strict
ProtectHome=read-only

ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes

LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictRealtime=yes
SystemCallArchitectures=native

RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
ReadWritePaths=/var/lib/vps-sentry /run/vps-sentry

LogRateLimitIntervalSec=30s
LogRateLimitBurst=200

TimeoutStartSec=45s
