# deploy/systemd/vps-sentry.timer.d/override.conf
# /etc/systemd/system/vps-sentry.timer.d/override.conf
[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

RandomizedDelaySec=30s
AccuracySec=1s
Persistent=true
