# systemd/vps-sentry.service.d/ship.conf
[Service]
ExecStartPost=-/usr/local/bin/vps-sentry-ship
