# deploy/systemd/vps-sentry-ship.service.d/zz-public-state-dir.conf
[Service]
# Legacy host overrides may tighten StateDirectoryMode; force the published
# dashboard files to remain traversable by the web app user.
StateDirectoryMode=0711
