# deploy/systemd/vps-sentry-web.service.d/40-garbage-cleanup-write-paths.conf
[Service]
# Allow the audited garbage-reclaim helper to remove only the safe cache/temp
# paths it targets, even though the web unit keeps home directories read-only.
ReadWritePaths=/tmp /var/tmp /var/cache/apt /var/log/journal /var/www -/root/.vscode-server -/root/.cache -/root/.npm -/root/.local/share/pnpm -/home/tony/.vscode-server -/home/tony/.cache -/home/tony/.npm -/home/tony/.local/share/pnpm -/home/tony/_backup/vps-sentry-web
