# deploy/systemd/vps-sentry.service.d/ship.conf
[Service]
ExecStartPost=
ExecStartPost=/bin/systemctl start --no-block vps-sentry-ship.service
